# brother-scand.service
[Unit]
Description=Brother Scanner Push-to-Scan Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=brother-scand
Group=brother-scand
ExecStart=/opt/brother-scand/brother_scan.py --config /etc/brother-scand/profiles.conf
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=10

# Output directory (override in profiles.conf or create this)
WorkingDirectory=/var/lib/brother-scand

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/brother-scand
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
